# com.advit.emailmcp.outlook-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.advit.emailmcp.outlook-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/advitarora/.local/bin/uv</string>
        <string>run</string>
        <string>--frozen</string>
        <string>--directory</string>
        <string>/Users/advitarora/Projects/email-mcp</string>
        <string>python</string>
        <string>refresh_outlook.py</string>
    </array>
    <key>StartInterval</key>
    <integer>72000</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/advitarora/.config/email-mcp/refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/advitarora/.config/email-mcp/refresh.err</string>
</dict>
</plist>
